# frigate_ingress_proxy/rootfs/etc/services.d/watchdog/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Watchdog service started"

while true; do
    sleep 30
done
