# timekpr/rootfs/etc/services.d/timekpr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Sets up the users .ssh folder to be persistent
if ! bashio::fs.directory_exists /app/ssh; then
    mkdir -p /app/ssh \
        || bashio::exit.nok 'Failed to create a persistent ssh folder'

fi
chmod 700 /app/ssh \
    || bashio::exit.nok \
        'Failed setting permissions on persistent ssh folder'

if bashio::config.has_value 'ssh_private_key'; then
    bashio::log.info "store ssh_private_key"

    rm -f /app/ssh/timekpr_ui_key
    while read -r line; do
        echo "$line" >> /app/ssh/timekpr_ui_key
    done <<< "$(bashio::config 'ssh_private_key')"
#   bashio::config 'ssh_private_key' > /app/ssh/timekpr_ui_key


    chmod 600 /app/ssh/timekpr_ui_key
fi

git -C /app pull \
    || bashio::exit.nok \
        'Failed update application from github'

## Run your program
#exec /usr/bin/my_program
#python /app/app.py
gunicorn --bind 0.0.0.0:5001 hass_wrapper_gunicorn:app --chdir /app/ --timeout 60
